# photoframe.conf
[program:photoframe]
directory=/home/pi/fantastic-photo-frame
command=python photoframe.py
autostart=true
autorestart=true
environment=DISPLAY=:0.0
user=pi